# ops/backup/helios-backup.cron.example
15 2 */2 * * /path/to/repo/ops/backup/helios-backup.sh >> /var/log/helios-backup.log 2>&1
