# ops/backup/helios-backup.cron.example
15 2 * * * /usr/local/bin/helios-backup.sh >> /var/log/helios-backup.log 2>&1
